# templates/restic.timer.j2
[Unit]
Description=Perform restic backup [{{ name }}]

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
